# server-generic/autoshutdown/files/cron.d/daily-shutdown
# TKR: Shutdown the system.
50 1 * * * root systemctl poweroff
